# {{cookiecutter.repo_name}}/provisioning/playbooks/roles/elasticsearch/templates/elasticsearch.conf
{% raw %}
[program:elasticsearch]
command=/usr/share/elasticsearch/bin/elasticsearch -Des.http.port=9200 -Des.config=/etc/elasticsearch/elasticsearch.yml -Des.path.conf=/etc/elasticsearch -Des.path.home=/usr/share/elasticsearch -Des.path.logs=/var/log/elasticsearch -Des.path.data=/var/lib/elasticsearch -Des.path.work=/tmp/elasticsearch -p /var/run/elasticsearch.pid
numprocs=1
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/elasticsearch_stdout.log
stderr_logfile=/var/log/supervisor/elasticsearch_stderr.log
{% endraw %}
